# containers/postgresql/cron.daily/backup-movies-db
#!/bin/sh

DB=movies
DATE=`date +%Y%m%d`
FILE="${DB}-${DATE}.plain"

pg_dump -h 127.0.0.1 -U postgres -w -f /mnt/backups/${FILE} ${DB}
gzip /mnt/backups/${FILE}
