# containers/postgresql/cron.daily/cleanup-old-backups
#!/bin/sh

cd /mnt/backups
find . -mtime 7 -delete
